# local/templates/cron.d/oud
# Trivadis AG, Business Development & Support (BDS)
# Saegereistrasse 29, 8152 Glattbrugg, Switzerland
# -------------------------------------------------------------------------
# Name.......: oud
# Author.....: Stefan Oehrli (oes) stefan.oehrli@trivadis.com
# Editor.....: Stefan Oehrli
# Date.......: 2018.03.18
# Revision...: --
# Purpose....: Cronjob File for OUD Environment
# Notes......: --
# Reference..: https://github.com/oehrlis/oudbase
# License....: GPL-3.0+
# -----------------------------------------------------------------------
# Modified...:
# see git revision history with git log for more information on changes
# -----------------------------------------------------------------------
# -----------------------------------------------------------------------
# weekly full backup of all runing instances
# -------------------------------------------------------------------------
0 19 * * 1   USER OUD_LOCAL/bin/oud_backup.sh -t FULL -m oud@oradba.ch >/dev/null 2>&1

# -----------------------------------------------------------------------
# daily incremental backup of all runing instances
# -----------------------------------------------------------------------
0 19 * * 0,2,3,4,5,6   USER OUD_LOCAL/bin/oud_backup.sh -t INCREMENTAL -m oud@oradba.ch >/dev/null 2>&1

# -----------------------------------------------------------------------
# daily export of all runing instances
# -----------------------------------------------------------------------
0 22 * * *   USER OUD_LOCAL/bin/oud_export.sh -m oud@oradba.ch >/dev/null 2>&1
